# deploy/signaldesk.service
[Unit]
Description=Signal Desk trading terminal
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=/opt/signaldesk
ExecStart=/usr/bin/python3 /opt/signaldesk/xauusd_trader.py
# Secrets (SIGNALDESK_PASS, optional SIGNALDESK_USER) live here, not in the unit
EnvironmentFile=-/etc/signaldesk.env
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
